# templates/Asq.crontab
SHELL=/bin/bash
MAILTO=sys-msgs@yourdomain.com
*/1 * * * * asq curl http://localhost:3000/refresh >/dev/null 2>&1
